# bin/block/crontab
* * * * * bash /var/root/bin/block.sh check
